# src/com.local.memory-cache.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.memory-cache</string>
	<key>Program</key>
	<string>__HOME__/.local/bin/create_memory_cache.sh</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/memory-cache.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/memory-cache.err.log</string>
</dict>
</plist>
